# systemd/sysoracle.service
[Unit]
Description=SysOracle Monitor

[Service]
ExecStart=/usr/local/bin/sysoracle
Restart=always

[Install]
WantedBy=default.target
